# Files/etc/s6/openssh/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

sudo -u root pushd /root > /dev/null
    exec sudo -u root /usr/sbin/sshd -D -e 2>&1
sudo -u root popd
