# Files/etc/s6/gitea/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd /data/gitea > /dev/null
    exec /usr/local/bin/gitea web
popd
